# backend/pizza-ranch-manager-backend.service
[Unit]
Description=Pizza Ranch FunZone Manager Backend API
After=network.target

[Service]
Type=simple
User=clarence
WorkingDirectory=/home/clarence/pizza-ranch-manager/backend
ExecStart=/usr/bin/python3 -m uvicorn server:app --host 0.0.0.0 --port 8082
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
